# supervisord.d/application.ini
[program:application]
command=/root/data/minitrill/code/backend/venv/bin/python application.py
directory=/root/data/minitrill/code/backend
user=root
autostart = true
autorestart=true
redirect_stderr=true
stdout_logfile=/root/data/minitrill/code/backend/log/application.log
stdout_logfile_maxbytes=200MB
killasgroup=false
stopsignal=INT
