# supervisord.d/redispub.ini
[program:redispub]
command=/root/data/minitrill/code/backend/venv/bin/python redispub.py
directory=/root/data/minitrill/code/backend
user=root
autostart = true
autorestart=true
redirect_stderr=true
stdout_logfile=/root/data/minitrill/code/backend/log/redispub.log
stdout_logfile_maxbytes=200MB
killasgroup=false
stopsignal=INT
